# build_pkg/DEBIAN/prerm
#!/bin/bash
set -e
systemctl disable helloworldwebapp.service
systemctl stop helloworldwebapp.service
